# packaging/systemd/digilol-cert-pushpuller-push.timer
[Unit]
Description=Push certificates to S3 daily

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
