# packaging/scripts/deb/prerm
#!/bin/sh
set -e

if [ -d /run/systemd/system ]; then
  systemctl disable digilol-cert-pushpuller-push.timer digilol-cert-pushpuller-pull.timer >/dev/null 2>&1 || true
fi

exit 0
